# scripts/etc/cron.d/pooget-autoupdate
# Poodillion Auto-Update Cron Job
# Checks for package updates periodically
#
# Format: minute hour day month dayofweek command
# */5 = every 5 minutes

# Check for updates every 5 minutes
*/5 * * * * root /usr/sbin/pooget-autoupdate-cron

# Note: Different systems can have different intervals
# underground.bbs: */2 (every 2 minutes - aggressive)
# university.edu: */10 (every 10 minutes - cautious)
# Others: */5 (every 5 minutes - default)
