# distro/rootfs/etc/inittab
# PoodillionOS inittab
# Format: id:runlevels:action:process

# System initialization
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/mount -t devtmpfs devtmpfs /dev
::sysinit:/bin/hostname -F /etc/hostname

# Start getty on tty1
tty1::respawn:/bin/getty /dev/tty1

# Shutdown actions
::shutdown:/bin/sync
::shutdown:/bin/umount -a
